# reminder.service
[Unit]
Description=Reminderscript

[Service]
ExecStart={basepath}/run.sh --headless
Restart=on-failure
